# templates/bravo-ebs-discovery.service.j2
[Unit]
Before=bravo.mount
After=boot.mount
Description=Create EnvironmentFile for bravo.mount

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gen_ebs_mount_env.sh
